# proxy/matrix-portal-proxy.service
[Unit]
Description=Matrix Portal HTTP Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/matrix-portal-tracker/proxy
ExecStart=/usr/bin/python3 /home/pi/matrix-portal-tracker/proxy/server.py
Restart=always
RestartSec=5
Environment=PORT=6590

[Install]
WantedBy=multi-user.target
